# supervisord.conf
[supervisord]
nodaemon=true
childlogdir=/app/logs

[program:nodejs]
command=/usr/bin/nodejs /app/app.js
